# test/invalid/bad-ports.container
[Container]
Image=nginx
# QL080
PublishPort=99999:80
