# test/invalid/bad-hosts.container
[Container]
Image=nginx
# QL081
AddHost=myhost-without-ip
