# test/invalid/bad-enums.container
[Container]
Image=nginx
# QL040
ReadOnly=invalid-bool
