# test/invalid/bad-bytes.container
[Container]
Image=nginx
# QL082
Memory=512mk
