# test/invalid/conditional-ctr.container
# QL061
[Container]
Image=nginx
Group=1000
# Missing User when Group is set
