# test/invalid/conflicts.container
[Container]
Image=nginx
# QL070
Rootfs=/path/to/rootfs
