# test/invalid/conflicts-reload.container
[Container]
Image=nginx
ReloadCmd=kill -SIGHUP 1
# QL070
ReloadSignal=SIGHUP
