# test/invalid/duplicate-keys.container
[Container]
ContainerName=web-app
# QL020
ContainerName=web-app-backup
Image=nginx
