# test/invalid/syntax.container
# QL002
Image=someimage

[Container]
# QL001
invalid_line_without_equal_sign
