# test/invalid/unknowns.container
# QL010
[Continer]
Image=nginx

[Container]
# QL030
UnknownKeyOption=value
