# test/valid/redis-runner.container
[Unit]
Description=Redis Service Runner

[Container]
Image=redis-custom:latest
PublishPort=6379:6379
